# init/alpine/bridge
#!/sbin/openrc-run

dir="/etc/void"
command="${dir}/bridge"
command_background=true
pidfile="/var/run/bridge.pid"

stop() {
 pkill -P $(cat "${pidfile}")
}

stop_post() {
 "${dir}/nobridge"
}
