# init/alpine/dhcp
#!/sbin/openrc-run

. /etc/void/fun.sh || exit 1
. /etc/void/cfg.sh || exit 1

depend() {
 need bridge
}

start() {
 ebegin "starting dhcpcd"
 sleep_step="2"

 # wait for bridge dev
 while ! dev_exist "${br}"; do
 echo "dhcp: no ${br} dev exist"
 sleep "${sleep_step}"; done

 # and link on it
 while ip -br link show "${br}" | grep -q DOWN; do
 echo "dhcp: no link on ${br}"
 sleep "${sleep_step}"; done

 ip -6 addr flush dev "${br}" # fck ipv6

 dhcpcd -f /dev/null -b --ipv4only --noipv4ll \
--option domain_name_servers \
--script /usr/lib/dhcpcd/dhcpcd-run-hooks \
--nohook test,hostname,ntp.conf "${br}"

 eend $?
}

stop() {
 ebegin "killing dhcpcd"
  pkill dhcpcd || exit 0 # crash recovery
 eend $?
}

stop_post() {
 ebegin "flushing ip address"
  dev_exist "${br}" && ip -4 addr flush "${br}"
 eend $?
}
