# init/alpine/vcom0
#!/sbin/openrc-run
# permanent pty for nut upsd

command="/usr/bin/socat"
command_args="pty,link=/dev/${RC_SVCNAME}t,rawer pty,link=/dev/${RC_SVCNAME},rawer"
command_background=true
pidfile="/var/run/${RC_SVCNAME}"

stop() {
 ebegin "Stopping socat bridge /dev/${RC_SVCNAME} <-> /dev/${RC_SERVICE}t"
 pkill -F "${pidfile}" || exit 0
 eend $?
}

start_post() {
 /bin/chmod 660 -- $(readlink "/dev/${RC_SVCNAME}")
}
