# init/alpine/vcom0_remote
#!/sbin/openrc-run

remote="tcp:10.200.12.2:4196,forever"
svc="$(echo ${RC_SVCNAME} | cut -d'-' -f1 | cut -d'_' -f1)"

command="/usr/bin/socat"
command_args="-U open:/dev/${svc}t,raw,nonblock ${remote}"
command_background=true
pidfile="/var/run/${RC_SVCNAME}"

depend() {
 need ${svc}
}

stop() {
 ebegin "Stopping socat bridge /dev/${svc}t -> ${remote}"
 pkill -F "${pidfile}"
 eend $?
}
