# init/alpine/vcom0
#!/sbin/openrc-run
# permanent pty for nut upsd
command="/usr/bin/socat"
command_args="pty,link=/dev/${RC_SVCNAME},raw,crnl pty,link=/dev/${RC_SVCNAME}t,raw,crnl"
command_background=true
pidfile="/var/run/${RC_SVCNAME}"

stop() {
 ebegin "Stopping socat bridge /dev/${RC_SVCNAME} -> /dev/${RC_SVCNAME}t"
 pkill -F "${pidfile}"
 eend $?
}

start_post() {
 /bin/chmod 660 -- $(readlink "/dev/${RC_SVCNAME}")
}
